# s6-rc.d/web-viewer/run
#!/command/execlineb -P
foreground { sleep 10 }
python3 /opt/web_viewer.py
