# s6-rc.d/viam-server/run
#!/command/execlineb -P
if { test -f /etc/viam.json }
foreground { sleep 10 }
/usr/local/bin/viam-server -config /etc/viam.json
