# s6-rc.d/can-spawner/run
#!/command/execlineb -P
foreground { sleep 10 }
python3 /opt/can_spawner.py
